# UnionTech-repoinfo-1.0/service/repoinfo.timer
[Unit]
Description=repoinfo timer

[Timer]
OnBootSec=10min
OnUnitInactiveSec=24h
Unit=repoinfo.service

[Install]
WantedBy=multi-user.target

