# UnionTech-repoinfo-1.0/service/repoinfo.timer
[Unit]
Description=repoinfo timer

[Timer]
OnCalendar=*-*-* 00:59:00
Unit=repoinfo.service

[Install]
WantedBy=multi-user.target

